# deploy/vps/supervisor/celerybeat.conf
; Supervisor config for StrateKaz Celery Beat (Scheduler)
; /etc/supervisor/conf.d/stratekaz-celerybeat.conf

[program:stratekaz-celerybeat]
command=/var/www/stratekaz/venv/bin/celery -A config beat --loglevel=INFO --scheduler django_celery_beat.schedulers:DatabaseScheduler
directory=/var/www/stratekaz/repo/backend
user=stratekaz
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/stratekaz/logs/celerybeat.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
environment=
    DJANGO_SETTINGS_MODULE="config.settings",
    PYTHONPATH="/var/www/stratekaz/repo/backend"

; Only one beat instance should run
startsecs=10
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=10
