# deploy/vps/supervisor/celery.conf
; Supervisor config for StrateKaz Celery Worker
; /etc/supervisor/conf.d/stratekaz-celery.conf

[program:stratekaz-celery]
command=/var/www/stratekaz/venv/bin/celery -A config worker --loglevel=INFO --concurrency=2
directory=/var/www/stratekaz/repo/backend
user=stratekaz
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/stratekaz/logs/celery.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=
    DJANGO_SETTINGS_MODULE="config.settings",
    PYTHONPATH="/var/www/stratekaz/repo/backend"

; Graceful stop
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=60

; Restart on memory limit (optional)
; startsecs=10
