# deploy/vps/supervisor/gunicorn.conf
; Supervisor config for StrateKaz Gunicorn
; /etc/supervisor/conf.d/stratekaz-gunicorn.conf

[program:stratekaz-gunicorn]
command=/var/www/stratekaz/venv/bin/gunicorn config.wsgi:application -c /var/www/stratekaz/gunicorn.conf.py
directory=/var/www/stratekaz/repo/backend
user=stratekaz
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/stratekaz/logs/gunicorn-supervisor.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=
    DJANGO_SETTINGS_MODULE="config.settings",
    PYTHONPATH="/var/www/stratekaz/repo/backend"

; Graceful stop
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=10
